# VM отдельно/postgres/roles/database/files/borg-backup.timer
[Unit]
Description=Automated Borg Backup Timer

[Timer]
OnCalendar=*-*-* 00:35

[Install]
WantedBy=timers.target
